# php/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
redirect_stderr=true
startretries=30
startsecs=10
