# php/supervisor/conf.d/php.conf
[program:php]
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true
redirect_stderr=true
startretries=30
startsecs=10
