# php/supervisor/conf.d/websocket.conf
[program:websocket]
command=/usr/local/bin/php /var/www/html/ctc/websocket/start.php -d
autostart=true
autorestart=true
redirect_stderr=true
