# php/cron.d/course
* * * * * www-data /usr/local/bin/php /var/www/html/ctc/scheduler.php 1>/dev/null 2>&1
