# php/cron.d/course
* * * * * root  echo "123456" >> /tmp/number
